# systemd/tgbot.service
# Example configuration for running the bot on a server as a systemd service
[Unit]
Description=YoutubeMusicDownloadBot daemon
After=network.target

[Service]
User=tgbot
Group=tgbot
Type=simple
WorkingDirectory=/opt/YoutubeMusicDownloadBot
ExecStart=/opt/YoutubeMusicDownloadBot/venv/bin/python bot.py
Restart=always

[Install]
WantedBy=multi-user.target
